# etc/crontab
#update GeoIP-db every month
0 1 2 * * /home/radio/bin/getgeoip.sh
*/5 * * * * /home/radio/bin/genGraph.pl 2> /dev/null
