# etc/crontab
#update GeoIP-db every month
0 1 2 * * /srv/radio/bin/getgeoip.sh
